# etc/cron.d/wakealarm
* * * * *   root    test -x /usr/sbin/wakealarm && /usr/sbin/wakealarm
